# util/etc/update-motd.d/50-greeting
#!/bin/sh
printf "%s \n" "  ___  ___ _  _ _____   __  _      _   ___      "
printf "%s \n" " / _ \| _ \ \| | __\ \ / / | |    /_\ | _ )     "
printf "%s \n" "| (_) |  _/ .` | _| \ V /  | |__ / _ \| _ \     "
printf "%s \n" " \___/|_| |_|\_|_|   \_/   |____/_/ \_\___/     "
printf "%s \n" "     "

